# cat-db.service
[Unit]
Name=Cat-DB
Description=a place for your cat images
After=network.target

[Service]
Type=simple
User=node
ExecStart=/usr/bin/node /home/node/cat-database/index.js 2>&1 /var/log/cat-db.log
Restart=on-failure

[Install]
WantedBy=multi-user.target